# extras/osdial.init
#!/bin/bash
#
# osdial  Startup script for OSDial services.
#
# chkconfig: - 64 36
# description:  OSDial services.
# processname:  osdial
# pidfile: /var/run/osdial/osdial.pid

