# extras/osdial.init
#!/bin/bash
#
# astguiclient  Startup script for VICIdial services.
#
# chkconfig: - 64 36
# description:  astGUIclient/VICIdial services.
# processname:  astguiclient
# pidfile: /var/run/astguiclient/astguiclient.pid

